# services/captive-watchdog.service
[Unit]
Description=Atlantis Captive Portal Watchdog
After=network-online.target wan-manager.service

[Service]
ExecStart=/usr/local/bin/captive-watchdog.sh
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=captive-watchdog

[Install]
WantedBy=multi-user.target
