# services/wan-manager.service
[Unit]
Description=Atlantis Travel Router WAN/LAN Manager
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wan-manager.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
